# devel_isolated/rosbash/etc/catkin/profile.d/15.rosbash.bash
# generated from rosbash/env-hooks/15.rosbash.bash.em

. "/home/drocode/ros_catkin_ws/src/ros/rosbash/rosbash"
