# devel/etc/catkin/profile.d/15.rosbash.fish
# generated from rosbash/env-hooks/15.rosbash.fish.em

. "/home/drocode/ros_catkin_ws/src/ros/rosbash/rosfish"
